# program-manager-testing.socket
[Socket]
ListenStream=/run/program-manager-testing.sk

[Install]
WantedBy=sockets.target